# scripts/hermes.service.tmpl
[Unit]
Description=Nous Research Hermes Agent (unattended gateway)
Documentation=https://github.com/nousresearch/hermes-agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=${HERMES_USER}
Group=${HERMES_USER}
WorkingDirectory=${WORKSPACE_DIR}

# Runtime env: Gemini key + Hermes locations. Root-owned, mode 0600.
EnvironmentFile=/etc/hermes-agent/hermes.env
Environment=HERMES_HOME=${HERMES_CONFIG_DIR}
Environment=XDG_STATE_HOME=${HERMES_CONFIG_DIR}/state
Environment=HERMES_GATEWAY_LOCK_DIR=${HERMES_CONFIG_DIR}/state/hermes/gateway-locks
# Keep the Hermes Python environment first so bundled skills that invoke
# `python` use the same interpreter that owns their optional dependencies.
# Both paths are included because upstream installations have used venv/ and
# .venv/ layouts. Non-existent PATH entries are harmless.
Environment=PATH=${HERMES_CONFIG_DIR}/hermes-agent/venv/bin:${HERMES_CONFIG_DIR}/hermes-agent/.venv/bin:${HERMES_HOME}/.local/bin:/usr/local/bin:/usr/bin:/bin

# Unattended background entrypoint. `hermes gateway` (no subcommand) runs the
# messaging gateway in the foreground -- there is no `gateway run` subcommand.
# It activates any platform (e.g. Telegram) purely from env vars being present
# (TELEGRAM_BOT_TOKEN, etc., via the EnvironmentFile below) and works with no
# platform configured too (standalone, bound to localhost) -- suitable for a
# Type=simple service. It does NOT open a public listener.
ExecStart=${HERMES_HOME}/.local/bin/hermes gateway

# Conservative restart policy for a small Free-Tier VM.
Restart=on-failure
RestartSec=10

# Log to journald (journald size caps are configured in install-service.sh).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-agent

# Light hardening.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=${HERMES_CONFIG_DIR} ${WORKSPACE_DIR}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
